# systemd/log-wagman.timer
[Unit]
Description=Logs Wagman status periodically.

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
