# systemd/waggle-log-wagman.timer
[Unit]
Description=Logs Wagman status periodically.

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
