# nc-wag-os/waggled/supervisor/waggle_reverse_ssh_tunnel.conf
[program:waggle_reverse_ssh_tunnel]
command=autossh.sh
directory=/usr/lib/waggle/nodecontroller/scripts/
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/waggle/supervisor_reverse_ssh_tunnel.err.log
stdout_logfile=/var/log/waggle/supervisor_reverse_ssh_tunnel.out.log
user=waggle
stopasgroup=true
killasgroup=true