# nc-wag-os/waggled/upstart/waggle_reverse_ssh_tunnel.conf
start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel[016]
respawn
exec /usr/lib/waggle/nodecontroller/scripts/autossh.sh