# etc/systemd/system/waggle-status.service
[Unit]
Description=Waggle status

[Service]
ExecStart=/usr/bin/python3 -u /usr/lib/waggle/core/scripts/status-service --beehive --silent
