# etc/systemd/system/waggle-cleanup-system.service
[Unit]
Description=Automatically cleanup system files

[Service]
ExecStart=/usr/lib/waggle/nodecontroller/scripts/waggle-cleanup-system
