# etc/systemd/system/waggle-status.service
[Unit]
Description=Waggle status

[Service]
ExecStart=/usr/lib/waggle/core/scripts/status-service --beehive