# etc/network/if-up.d/firewall
#!/bin/bash

iptables-restore < /etc/firewall.rules
