# upstart/waggle_reverse_ssh_tunnel.conf
start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel[016]
respawn
respawn limit unlimited
post-stop exec sleep 5
exec /usr/lib/waggle/nodecontroller/scripts/autossh.sh