# upstart/waggle_data_cache.conf

start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel[016]
respawn
respawn limit unlimited
pre-start script
  set +e
  for pid in $(ps auxwww | grep 'python Data_Cache.py' | grep -oP "^\S*\s*\d*" | grep -oP "\d*$") ; do echo "Killing exsiting process: kill -9 ${pid}" ; kill -9 $pid ; done
end script
post-stop exec sleep 5
chdir /usr/lib/waggle/nodecontroller/nc-wag-os/waggled/DataCache
exec python Data_Cache.py --force --logging


